# conf/supervisor.conf.d/web.conf
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:uwsgi]
user=app
directory=/app/ecs
environment=HOME="/app",USER="app",VIRTUAL_ENV="/app/env",PATH="/app/env/bin:%(ENV_PATH)s"
command=/app/env/bin/uwsgi --yaml /app/ecs/conf/uwsgi.yml
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:uwsgi_exporter]
user=app
directory=/app/ecs
environment=HOME="/app",USER="app",VIRTUAL_ENV="/app/env",PATH="/app/env/bin:%(ENV_PATH)s"
command=/bin/bash -c 'sleep 10; exec /app/env/bin/uwsgi_exporter.py --stats_addr 127.0.0.1:8017 --bind_address "0.0.0.0" --bind_port 1717'
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
